# apps/docker-plex/root/etc/services.d/plexpatch/run
#!/usr/bin/with-contenv bash

echo "Starting Plex Media Server Patching."
umask 022
command="/app/plexpatch.sh"
sleep 2
exec s6-setuidgid abc $command
