# apps/docker-uploader/root/etc/services.d/uploader/run
#!/usr/bin/with-contenv bash
#####################################
# All rights reserved.              #
# started from Zero                 #
# Docker owned dockserver           #
# Docker Maintainer dockserver      #
#####################################
#####################################
# THIS DOCKER IS UNDER LICENSE      #
# NO CUSTOMIZING IS ALLOWED         #
# NO REBRANDING IS ALLOWED          #
# NO CODE MIRRORING IS ALLOWED      #
#####################################

umask 022

chown -cR abc:abc \
	/app \
	/system/*

chmod -cR 755 /app \
        /system/*

command="/app/uploader/uploader.sh"
echo "-> Starting Multi Thread Uploader <-"
exec s6-setuidgid abc $command
